# moon_profile_runner/packaging/moon-profile-runner-git.install
post_install() {
	echo ":: Launch MoonProfile Runner once (applications menu, or run 'moon-profile-runner')"
	echo "   and it registers itself for autostart on every future login."
	echo "   Or do it by hand now: systemctl --user enable --now moon-profile-runner.service"
}

post_upgrade() {
	post_install
}
